# wk-enhanced-api/deploy/wk-enhanced-api.service
# systemd unit that brings up the wk-enhanced-api Docker Compose stack on boot.
#
# Install at: /etc/systemd/system/wk-enhanced-api.service
# Then:       systemctl daemon-reload && systemctl enable --now wk-enhanced-api
# Watch:      docker compose -f /opt/wk-enhanced-api/wk-enhanced-api/compose.yaml logs -f api
#         or: journalctl -fu wk-enhanced-api    (for unit-level start/stop events)
#
# Container-side concerns (the unprivileged `bun` user, sandbox flags,
# binary location) are handled by the Dockerfile + Docker daemon. This
# unit's only job is to invoke `docker compose up` once on boot, and
# `docker compose down` cleanly on shutdown. Container auto-restart is
# governed by `restart: unless-stopped` in compose.yaml, not by systemd.
#
# Pre-Docker droplets: see deploy/README.md "Migrating from a pre-Docker
# droplet" for the one-time conversion steps (install Docker, chown
# /var/lib, etc.). After conversion this same unit file runs and the
# old bare-metal Bun + chunk of Protect*= hardening goes away.

[Unit]
Description=wk-enhanced-api — backing API for the WKEnhanced userscript (Docker Compose)
Documentation=https://github.com/<your-user-or-org>/WaniKani/blob/main/wk-enhanced-api/README.md
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/wk-enhanced-api/wk-enhanced-api

# Point compose.yaml's parameterized inputs at their prod-only values.
# Locally these default to `./.env` and `./.compose-data` respectively,
# so `docker compose up` from a fresh checkout on a dev machine works
# without needing this override. On the droplet:
#   - ENV_FILE   → the chmod-600 root-only env file under /etc.
#   - DATA_DIR   → the persistent /var/lib path the host bind-mounts
#                  into the container; matches the container-side path
#                  so DATABASE_FILE in the env file resolves to the
#                  same bytes whether read via `docker exec` or off the
#                  host filesystem directly.
Environment="ENV_FILE=/etc/wk-enhanced-api/env"
Environment="DATA_DIR=/var/lib/wk-enhanced-api"

# Compose v2 is invoked as a `docker` sub-command (no hyphen). Standard
# install location varies between distros — `/usr/bin/docker` is the
# Debian/Ubuntu apt path; on hosts that installed via the official
# `get.docker.com` script the binary may live at `/usr/bin/docker` too,
# but check `which docker` if this path fails on a fresh provision.
ExecStart=/usr/bin/docker compose up -d --remove-orphans
ExecStop=/usr/bin/docker compose down

# Reload sends SIGHUP to the container if it knows how to handle it.
# Bun's HTTP server doesn't currently react to SIGHUP, so a config
# change still requires `systemctl restart wk-enhanced-api` (which runs
# down + up in sequence and gives us a fresh container picking up the
# updated env file). Documented here so future-us doesn't expect a
# magic graceful reload.
ExecReload=/usr/bin/docker compose up -d --force-recreate

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
