# wk-enhanced-api/deploy/wk-enhanced-api.service
# systemd unit for the wk-enhanced-api deployment.
#
# Install at: /etc/systemd/system/wk-enhanced-api.service
# Then:       systemctl daemon-reload && systemctl enable --now wk-enhanced-api
# Watch:      journalctl -fu wk-enhanced-api
#
# Runs as a dedicated unprivileged `wkenhanced` user. Create with:
#   useradd --system --no-create-home --shell /usr/sbin/nologin wkenhanced
#   install -d -o wkenhanced -g wkenhanced /var/lib/wk-enhanced-api
#   chown -R wkenhanced:wkenhanced /opt/wk-enhanced-api
#
# The unit invokes `bun run start` from package.json (which runs
# `bun run src/index.ts`). Going through the script rather than calling the
# entrypoint directly means an entrypoint rename in the future doesn't break
# the unit.
#
# Note: serves the source tree at wk-enhanced-api/ inside the cloned repo.
# WorkingDirectory points at /opt/wk-enhanced-api/wk-enhanced-api accordingly.
# (Production deployments that predate the 2026-05-25 rename had this at
# /opt/wk-enhanced-api/wk-vocab-api — `mv` that directory before
# `systemctl daemon-reload && systemctl restart wk-enhanced-api`.)

[Unit]
Description=wk-enhanced-api — backing API for the WKEnhanced userscript
Documentation=https://github.com/<your-user-or-org>/WaniKani/blob/main/wk-enhanced-api/README.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=wkenhanced
Group=wkenhanced
WorkingDirectory=/opt/wk-enhanced-api/wk-enhanced-api
EnvironmentFile=/etc/wk-enhanced-api/env

# /usr/local/bin/bun, NOT /root/.bun/bin/bun. The service runs as the
# unprivileged `wkenhanced` user, and our `ProtectHome=true` hardening below
# blocks reads of /root anyway — so Bun must live somewhere world-readable
# outside /home and /root. The deploy README's install step copies the
# official-installer-placed binary from /root/.bun/bin/bun into
# /usr/local/bin/bun for this reason.
ExecStart=/usr/local/bin/bun run start

Restart=always
RestartSec=5

# Logging: journald captures stdout/stderr; structured JSON lines stay
# greppable via `journalctl -fu wk-enhanced-api`.
StandardOutput=journal
StandardError=journal

# Hardening: the service only needs outbound HTTPS, the SQLite file under
# /var/lib/wk-enhanced-api, and read access to /opt/wk-enhanced-api. Lock
# down everything else.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/wk-enhanced-api
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=false  # Bun's JIT needs W+X memory.

[Install]
WantedBy=multi-user.target
