# wk-enhanced-api/deploy/wk-enhanced-api-backup.service
# Daily SQLite backup. Runs deploy/backup.ts INSIDE the running container
# via `docker exec` — the script needs bun:sqlite (for VACUUM INTO) and
# Bun.S3Client (for the upload + retention pruning), both of which live
# in the container image, not on the host.
#
# Install at: /etc/systemd/system/wk-enhanced-api-backup.service
# Paired with: wk-enhanced-api-backup.timer (the schedule).
# Manual run:  systemctl start wk-enhanced-api-backup
# Watch:       journalctl -fu wk-enhanced-api-backup
#
# The script picks up DATABASE_FILE + the S3_* creds from the container's
# env (sourced from /etc/wk-enhanced-api/env via compose.yaml's env_file:),
# so this unit needs no EnvironmentFile= of its own.
#
# Failure surfaces in `systemctl list-units --failed`. The backup script
# exits non-zero on any error → docker exec exits non-zero → this oneshot
# fails and stays visible until you run `systemctl reset-failed`.

[Unit]
Description=Daily SQLite backup of wk-enhanced-api to DO Spaces
Documentation=https://github.com/<your-user-or-org>/WaniKani/blob/main/wk-enhanced-api/README.md
Requires=wk-enhanced-api.service
After=wk-enhanced-api.service

[Service]
Type=oneshot

# The container is named via compose.yaml's `container_name: wk-enhanced-api`.
# `docker exec` runs the script as the container's USER (uid 1000), which
# already has the right perms on the bind-mounted /var/lib/wk-enhanced-api.
# No additional User= / hardening flags needed — the container provides
# the sandbox.
ExecStart=/usr/bin/docker exec wk-enhanced-api bun run /app/deploy/backup.ts

StandardOutput=journal
StandardError=journal
